# .deploy/supervisord.conf
[program:uwsgi-app]
command = /usr/local/bin/uwsgi --ini /code/.deploy/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx